# stow/khalid/.config/systemd/user/checkmail.timer
[Unit]
Description=Check Mail every five minutes
RefuseManualStart=no
RefuseManualStop=no

[Timer]
Persistent=false
OnStartupSec=30sec
OnUnitActiveSec=10min
Unit=checkmail.service

[Install]
WantedBy=default.target