# khalid/.config/systemd/user/checkbatt.service
[Unit]
Description=check battery
RefuseManualStart=no
RefuseManualStop=yes

[Service]
Type=oneshot
ExecStart=%h/bin/checkbatt.sh